# amavis-stats.cron
#
# cron job for the amavis-stats package
#
*/5 * * * *	root	/usr/bin/amavis-stats /var/log/maillog
