# amavis-stats.cron
#
# cron job for the amavis-stats package
#
*/5 * * * *	root	LC_ALL="C" /usr/sbin/amavis-stats /var/log/maillog
